# templates/com.claude.ai-research-assistant.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.claude.ai-research-assistant</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/jens.echterling/GitHub/Development/AI research assistant/scripts/run.sh</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>6</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/jens.echterling/.claude/logs/ai-research-assistant.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/jens.echterling/.claude/logs/ai-research-assistant.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
        <key>HOME</key>
        <string>/Users/jens.echterling</string>
    </dict>
</dict>
</plist>
